# usbMitm.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/ubuntu/usb-mitm/connection.bash

[Install]
WantedBy=multi-user.target